# ts/5.1/packages/hzdr/etc/udev/rules.d/01-usbsound.rules
ACTION=="add", SUBSYSTEM=="usb", DRIVERS=="snd-usb-audio", RUN+="/bin/alsaconfigure.sh"   
ACTION=="remove", SUBSYSTEM=="sound", ENV{ID_BUS}=="usb", RUN+="/bin/alsaconfigure.sh"

